# apps/core/fazle-core.service
[Unit]
Description=Fazle Core — WhatsApp AI Backend
After=network.target docker.service

[Service]
Type=simple
User=azim
WorkingDirectory=/home/azim/core
EnvironmentFile=/home/azim/core/.env
EnvironmentFile=-/home/azim/secure-env-backup/runtime-services.env
ExecStartPre=-/home/azim/shared/scripts/generate_runtime_services_env.sh
ExecStart=/home/azim/.venv/bin/python run.py
Restart=always
RestartSec=5
StandardOutput=append:/home/azim/core/logs/fazle-core.log
StandardError=append:/home/azim/core/logs/fazle-core-error.log

[Install]
WantedBy=multi-user.target
